# @ripper/.config/systemd/user/offlineimap-oneshot@tamu.service.d/vault.conf
[Unit]
Requires=mnt-vault.mount
After=mnt-vault.mount
